# v2ray/tls-ws/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:v2ray]
command=v2ray -config=/etc/v2ray/config.json

[program:nginx]
command=nginx -g "daemon off;"